# build/COPY_ROOT_PROTON/etc/supervisor/supervisord.d/proton-services.conf
[program:proton-init]
command=/opt/ai-dock/bin/proton-init.sh
directory=/opt/ai-dock/bin
autostart=true
autorestart=false
priority=900
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/proton-init.log
stderr_logfile=/var/log/supervisor/proton-init.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile_backups=1
user=%(ENV_WORKSPACE_USER)s

[program:gamemode]
command=/usr/bin/gamemoded -F
directory=/
autostart=%(ENV_ENABLE_GAMEMODE)s
autorestart=true
priority=950
stdout_logfile=/var/log/supervisor/gamemode.log
stderr_logfile=/var/log/supervisor/gamemode.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile_backups=1
user=%(ENV_WORKSPACE_USER)s
environment=HOME="%(ENV_WORKSPACE)s"